# wave-sample/rootfs/etc/services.d/nginx/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: Wave Engine Sample
# Runs the Nginx daemon
# ==============================================================================

bashio::log.info "Starting NGinx..."
exec nginx